# service/redis/run
#!/bin/sh
exec chpst -u dvv /usr/local/bin/redis-server redis.conf 2>&1
