# runit/webify/run
#!/bin/sh
exec webfsd -F -r /home/web 2>&1
