# service/haproxy/run
#!/bin/sh
exec /usr/local/bin/haproxy -f haproxy.conf 2>&1
