# runit/flashpolicy/run
#!/bin/sh
exec perl -T socketpolicy.pl 2>&1
