# service/stud/run
#!/bin/sh

ulimit -n 100000
exec stud -n 1 -u stud -q dummy -f \*,443 -b 127.0.0.1,65443 --write-proxy -B 2048 -C 10000 -c ALL -- farm.pem 2>&1
