# runit/stud/run
#!/bin/sh
exec /usr/local/bin/stud -n 1 -u stud -f \*,443 -b 127.0.0.1,65443 --write-proxy -B 100 -C 100 farm.pem 2>&1
